# deploy/trading_bot.service
[Unit]
Description=Autonomous Trading Bot Daemon
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{APP_DIR}}
ExecStart={{APP_DIR}}/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
